# etc/systemd/system/oci-secretsync.service
[Unit]
Description=Sync secrets from OCI Vault to /run/secrets
Wants=network-online.target
After=network-online.target run-secrets.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/oci-secretsync /etc/oci-secretsync/map.csv
User=root
Group=root
# noisy failures shouldn't reboot the box:
SuccessExitStatus=0 3

[Install]
WantedBy=multi-user.target
